# deploy/pockethost-panel.service.example
[Unit]
Description=PocketHost management panel (tailnet-only)
After=network-online.target pockethost.service
Wants=network-online.target

[Service]
Type=simple
User=user
Group=user
WorkingDirectory=/home/user/services/pockethost-panel
Environment=HOME=/home/user
Environment=PATH=/home/user/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PANEL_PORT=8096
Environment=PH_HOME=/home/user/services/pockethost/home
Environment=PH_PORT=8095
Environment=PH_APEX_DOMAIN=pocket.example.com
Environment=PH_SERVICE=pockethost
Environment=PH_AUTO_ADMIN_LOGIN=admin@example.com
# Secrets — fill in before deploying. Do NOT commit the real values.
Environment=PH_AUTO_ADMIN_PASSWORD=__SET_ME__
Environment=PANEL_USER=admin
Environment=PANEL_PASS=__SET_ME__
ExecStart=/home/user/.bun/bin/bun run src/server.js
Restart=always
RestartSec=5
StandardOutput=append:/home/user/services/pockethost-panel/panel.log
StandardError=append:/home/user/services/pockethost-panel/panel.log

[Install]
WantedBy=multi-user.target
